# infra/systemd/openclaw-gateway.service
[Unit]
Description=OpenClaw Gateway
Documentation=https://github.com/OpenCodeEngineer/webagent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=${APP_USER}
Group=${APP_USER}
WorkingDirectory=${APP_DIR}
EnvironmentFile=-${APP_DIR}/.env
ExecStart=/usr/bin/pnpm --filter @openclaw/gateway start

Restart=on-failure
RestartSec=5s
StartLimitInterval=60s
StartLimitBurst=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=openclaw-gateway

LimitNOFILE=65536
LimitNPROC=4096
MemoryMax=512M
CPUQuota=100%

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=${APP_DIR}

[Install]
WantedBy=multi-user.target
