# infra/systemd/webagent-admin.service
[Unit]
Description=WebAgent Admin UI (port 3000)
Documentation=https://github.com/OpenCodeEngineer/webagent
After=network-online.target webagent-proxy.service
Wants=network-online.target
Requires=webagent-proxy.service

[Service]
Type=simple
User=${APP_USER}
Group=${APP_USER}
WorkingDirectory=${APP_DIR}/packages/admin
EnvironmentFile=-${APP_DIR}/.env
Environment=PORT=3000
ExecStart=/bin/bash -lc 'set -euo pipefail; if [[ ! -d "${APP_DIR}/packages/admin/.next/static" ]]; then echo "webagent-admin missing static assets directory: ${APP_DIR}/packages/admin/.next/static" >&2; exit 1; fi; if [[ -f "${APP_DIR}/packages/admin/.next/standalone/server.js" ]]; then exec /usr/bin/env PORT=3000 /usr/bin/node "${APP_DIR}/packages/admin/.next/standalone/server.js"; fi; if [[ -f "${APP_DIR}/packages/admin/.next/standalone/packages/admin/server.js" ]]; then exec /usr/bin/env PORT=3000 /usr/bin/node "${APP_DIR}/packages/admin/.next/standalone/packages/admin/server.js"; fi; echo "webagent-admin missing Next standalone entrypoint. Expected one of:" >&2; echo "  ${APP_DIR}/packages/admin/.next/standalone/server.js" >&2; echo "  ${APP_DIR}/packages/admin/.next/standalone/packages/admin/server.js" >&2; echo "Run build first: pnpm build" >&2; exit 1'

Restart=on-failure
RestartSec=5s
StartLimitInterval=60s
StartLimitBurst=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=webagent-admin

LimitNOFILE=65536
LimitNPROC=4096
MemoryMax=512M
CPUQuota=100%

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=${APP_DIR}

[Install]
WantedBy=multi-user.target
