# infra/systemd/webagent-proxy.service
[Unit]
Description=WebAgent Proxy (port 3001)
Documentation=https://github.com/OpenCodeEngineer/webagent
After=network-online.target openclaw.service
Wants=network-online.target
Requires=openclaw.service

[Service]
Type=simple
User=${APP_USER}
Group=${APP_USER}
WorkingDirectory=${APP_DIR}
EnvironmentFile=-${APP_DIR}/.env
ExecStart=/usr/bin/pnpm --filter @webagent/proxy start

Restart=on-failure
RestartSec=5s
StartLimitInterval=60s
StartLimitBurst=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=webagent-proxy

# Default-safe bind host for backend service.
Environment=PROXY_BIND_HOST=127.0.0.1

LimitNOFILE=65536
LimitNPROC=4096
MemoryMax=1G
CPUQuota=200%

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=${APP_DIR}

[Install]
WantedBy=multi-user.target
