# infra/systemd/openclaw.service.d/override.conf
[Service]
EnvironmentFile=-${APP_DIR}/.env
Environment=OPENCLAW_CONFIG_PATH=${APP_DIR}/openclaw/config/openclaw.json5
