# infra/systemd/openclaw.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/node /usr/local/lib/oc-gw.js gateway --port 18789
Environment=OPENCLAW_GATEWAY_TOKEN=${OPENCLAW_GATEWAY_TOKEN}
Environment=OPENCLAW_CONFIG_PATH=/opt/webagent/openclaw/config/openclaw.json5
Environment=AZURE_DEV_AI_API_KEY=${AZURE_DEV_AI_API_KEY}
Environment=AZURE_DEV_AI_BASE_URL=${AZURE_DEV_AI_BASE_URL}
WorkingDirectory=/home/openclaw
